# sources/amqp2engines/etc/engines/cleaner-alerts.ini
[program:engine-cleaner-alerts]

autostart=false

directory=%(ENV_HOME)s
numprocs=1
process_name=%(program_name)s-%(process_num)d

command=engine-launcher -e cleaner -n cleaner_alerts -w %(process_num)d -l info

stdout_logfile=%(ENV_HOME)s/var/log/engines/cleaner_alerts.log
stderr_logfile=%(ENV_HOME)s/var/log/engines/cleaner_alerts.log