# sources/python/canopsis/etc/engines/alerts.ini
[program:engine-alerts]

autostart=false

directory=%(ENV_HOME)s
numprocs=1
process_name=%(program_name)s-%(process_num)d

command=engine-launcher -e dynamic -n alerts -w %(process_num)d -l info

stdout_logfile=%(ENV_HOME)s/var/log/engines/alerts.log
stderr_logfile=%(ENV_HOME)s/var/log/engines/alerts.log