# sources/python/engines/etc/engines/crecord-dispatcher.ini
[program:engine-crecord-dispatcher]

autostart=false

directory=%(ENV_HOME)s
numprocs=1
process_name=%(program_name)s-%(process_num)d

command=engine-launcher -e crecord_dispatcher -n crecord_dispatcher -w %(process_num)d -l info

stdout_logfile=%(ENV_HOME)s/var/log/engines/crecord_dispatcher.log
stderr_logfile=%(ENV_HOME)s/var/log/engines/crecord_dispatcher.log