# sources/python/engines/etc/engines/event-filter.ini
[program:engine-event-filter]

autostart=false

directory=%(ENV_HOME)s
numprocs=1
process_name=%(program_name)s-%(process_num)d

command=engine-launcher -e event_filter -n event_filter -w %(process_num)d -l debug

stdout_logfile=%(ENV_HOME)s/var/log/engines/event_filter.log
stderr_logfile=%(ENV_HOME)s/var/log/engines/event_filter.log