# sources/python/engines/etc/engines/selector.ini
[program:engine-watcher]

autostart=false

directory=%(ENV_HOME)s
numprocs=1
process_name=%(program_name)s-%(process_num)d

command=engine-launcher -e dynamic -n watcher -w %(process_num)d -l info

stdout_logfile=%(ENV_HOME)s/var/log/engines/watcher.log
stderr_logfile=%(ENV_HOME)s/var/log/engines/watcher.log
