# sources/python/engines/etc/engines/ticket.ini
[program:engine-ticket]

autostart=false

directory=%(ENV_HOME)s
numprocs=1
process_name=%(program_name)s-%(process_num)d

command=engine-launcher -e ticket -n ticket -w %(process_num)d -l info

stdout_logfile=%(ENV_HOME)s/var/log/engines/ticket.log
stderr_logfile=%(ENV_HOME)s/var/log/engines/ticket.log