# community/sources/amqp2engines/etc/engines/alertcounter.ini
[program:engine-alertcounter]

autostart=false

directory=%(ENV_HOME)s
numprocs=1
process_name=%(program_name)s-%(process_num)d

command=engine-launcher -e alertcounter -n alertcounter -w %(process_num)d -l info

stdout_logfile=%(ENV_HOME)s/var/log/engines/alertcounter.log
stderr_logfile=%(ENV_HOME)s/var/log/engines/alertcounter.log