# sources/canolibs/engines/etc/engines/scheduler.ini
[program:engine-scheduler]

autostart=false

directory=%(ENV_HOME)s
numprocs=1
process_name=%(program_name)s-%(process_num)d

command=engine-launcher -e scheduler -n scheduler -w %(process_num)d -l info

stdout_logfile=%(ENV_HOME)s/var/log/engines/scheduler.log
stderr_logfile=%(ENV_HOME)s/var/log/engines/scheduler.log