# sources/canolibs/engines/etc/engines/taskhandler-mail.ini
[program:taskhandler-mail]

autostart=false

directory=%(ENV_HOME)s
numprocs=1
process_name=%(program_name)s-%(process_num)d

command=engine-launcher -e task_mail -n task_mail -w %(process_num)d -l info

stdout_logfile=%(ENV_HOME)s/var/log/engines/task_mail.log
stderr_logfile=%(ENV_HOME)s/var/log/engines/task_mail.log